# etc/profile.d/cbrief.csh
if ( $?HOME ) then
	setenv BACKUPDIR "${HOME}/.backup"
	setenv JED_HOME "${HOME}/.jed"
endif
