# etc/profile.d/cbrief.sh
if [ -n "$HOME" ]; then
	export BACKUPDIR=${HOME}/.backup
	export JED_HOME=${HOME}/.jed
fi
